# home/private_Library/private_LaunchAgents/com.rajiv.chezmoi-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rajiv.chezmoi-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>{{ .chezmoi.homeDir }}/bin/chezmoi-sync-check</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.cache/chezmoi-sync/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.cache/chezmoi-sync/daemon.log</string>
</dict>
</plist>
